# merecat.service.in
[Unit]
Description=Merecat web server
Documentation=man:merecat(8) man:ssi(8) man:redirect(8) man:htpasswd(1)
Documentation=file:@DOCDIR@/merecat/README.md
ConditionPathExists=/var/www
After=network-online.target
Requires=network-online.target

[Service]
Type=simple
ExecStart=@SBINDIR@/merecat -n -c '**.cgi' -r /var/www

[Install]
WantedBy=multi-user.target
